# SOURCES/drefd-dash-1.0/lib/systemd/system/drefd-dash.service
[Unit]
Description=drefd dashboard
After=network.target

[Service]
User=root
ExecStart=/opt/drefd-dash/drefd_dash.sh
Restart=always

[Install]
WantedBy=multi-user.target
